# deploy/com.binance.bot.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.binance.bot.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kh/binance/.venv/bin/python</string>
        <string>/Users/kh/binance/watchdog.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kh/binance</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kh/binance/logs/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kh/binance/logs/launchd.stderr.log</string>
</dict>
</plist>